# templates/dnf-automatic.timer.j2
[Timer]
OnCalendar={{ dnf_automatic_timer_oncalendar }}
RandomizedDelaySec={{ dnf_automatic_timer_randomized_delay }}
